# ihisengine/odb_cluster/supervisord.conf
[supervisord]
nodaemon = true
logfile = /opt/zato/supervisord.log
loglevel = debug
user = root
directory = /opt/zato

[program:zato-odb-cluster]
command = /bin/sh
user = zato
